# package/rooter/0routerspecfic/x3000/files/etc/hotplug.d/iface/20-wanup
#!/bin/sh

logger -t URL-DEBUG "hotplug (iface): action='$ACTION' interface='$INTERFACE'"

if [ "$ACTION" = ifdown -a "$INTERFACE" = "wan" ]; then
	sleep 5
	ifup wan
fi